# supervisor/163.gs-redis.conf
[program:163gs_redis]
user = projects
umask = 077
directory = /tmp
command = /usr/bin/redis-server /etc/redis/redis-163gs.conf
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/supervisor/163gs-redis.log
stderr_logfile = /var/log/supervisor/163gs-redis.error.log
